# luci-app-wififailover/root/etc/init.d/wififailover
#!/bin/sh /etc/rc.common

START=99
SCRIPT_FILE_NAME="wififailover"

start() {
    eval "/usr/libexec/wififailover/$SCRIPT_FILE_NAME &"
}

stop() {
    killall $SCRIPT_FILE_NAME
}