# luci-app-wififailover/files/etc/uci-defaults/99-wififailover
#!/bin/sh
# Файл: files/etc/uci-defaults/99-wififailover
# UCI инициализация для feeds пакета

# Проверка существования конфигурации
if ! uci -q get wififailover.settings >/dev/null 2>&1; then
    # Создание базовых настроек
    uci -q batch <<EOF
set wififailover.settings=global
set wififailover.settings.enabled='0'
set wififailover.settings.ping_host='8.8.8.8'
set wififailover.settings.check_interval='30'
set wififailover.settings.max_failures='3' 
set wififailover.settings.connection_timeout='60'
set wififailover.settings.debug='0'
set wififailover.settings.current_network=''
set wififailover.settings.auto_switch='1'
commit wififailover
EOF
    
    logger -t wififailover "UCI configuration initialized"
fi

# Добавление в LuCI меню если отсутствует
if ! uci -q get luci.wififailover >/dev/null 2>&1; then
    uci -q batch <<EOF
set luci.wififailover=internal
set luci.wififailover.lang='auto'
set luci.wififailover.mediaurlbase='/luci-static/resources'
set luci.wififailover.resourcebase='/luci-static/resources'
commit luci
EOF
fi

# Создание символической ссылки для совместимости
[ ! -d /www/luci-static ] && mkdir -p /www/luci-static
[ ! -L /www/luci-static/wififailover ] && \
    ln -sf /www/luci-static/resources /www/luci-static/wififailover 2>/dev/null || true

exit 0